# deploy/systemd/qyh-data-plane.service
[Unit]
Description=QYH Robot Data Plane (WebSocket/ROS2)
Documentation=https://github.com/qyh-robot/backend
After=network.target network-online.target
# ROS2 daemon is handled by middleware, not a system service usually.
# If you have a dedicated DDS discovery server, add it here.

[Service]
Type=simple
User=robot
Group=robot
WorkingDirectory=/opt/qyh-robot/data_plane
EnvironmentFile=/etc/qyh-robot/environment

# 使用 entrypoint 加载 ROS2 环境
ExecStart=/opt/qyh-robot/ros_entrypoint.sh /opt/qyh-robot/data_plane/bin/data_plane_server /etc/qyh-robot/data_plane/config.yaml

Restart=always
RestartSec=3

# 实时调度（可选，需要权限）
# Nice=-10
# CPUSchedulingPolicy=fifo
# CPUSchedulingPriority=50

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/qyh-robot
PrivateTmp=true

# 资源限制
MemoryMax=256M
CPUQuota=200%

[Install]
WantedBy=multi-user.target
