# deploy/systemd/qyh-control-plane.service
[Unit]
Description=QYH Robot Control Plane (FastAPI)
Documentation=https://github.com/qyh-robot/backend
After=network.target
Wants=qyh-data-plane.service qyh-media-plane.service

[Service]
Type=simple
User=robot
Group=robot
WorkingDirectory=/opt/qyh-robot/control_plane
Environment=PYTHONPATH=/opt/qyh-robot/control_plane
EnvironmentFile=/etc/qyh-robot/environment

ExecStart=/opt/qyh-robot/control_plane/.venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-config /etc/qyh-robot/control_plane/log_config.yaml \
    --log-level info

Restart=always
RestartSec=5

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/qyh-robot /var/log/qyh-robot
PrivateTmp=true

# 资源限制
MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
