# deploy/systemd/qyh-media-plane.service
[Unit]
Description=QYH Robot Media Plane (WebRTC/GStreamer)
Documentation=https://github.com/qyh-robot/backend
After=network.target

[Service]
Type=simple
User=robot
Group=robot
WorkingDirectory=/opt/qyh-robot/media_plane
EnvironmentFile=/etc/qyh-robot/environment

# GStreamer 环境
Environment=GST_DEBUG=2
Environment=GST_PLUGIN_PATH=/usr/lib/aarch64-linux-gnu/gstreamer-1.0

ExecStart=/opt/qyh-robot/media_plane/bin/media_plane_server \
    /etc/qyh-robot/media_plane/config.yaml

Restart=always
RestartSec=3

# 视频设备访问
SupplementaryGroups=video

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/qyh-robot /dev/video0 /dev/video1
PrivateTmp=true

# 资源限制（媒体需要更多资源）
MemoryMax=512M
CPUQuota=300%

[Install]
WantedBy=multi-user.target
